# .docker/web/crontab
PATH=$PATH:/var/www/scripts
*/15 * * * * scheduler >> /tmp/schedule.log
