# .docker/web/crontab
PATH=$PATH:/var/www/bin
*/15 * * * * gazelle scheduler >> /tmp/schedule.log
