# .docker/web/crontab
PATH=$PATH:/var/www/bin
*/15 * * * * scheduler >> /tmp/schedule.log
